# scripts/index-watchdog.service
[Unit]
Description=Restart Index when its Docker health check fails
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/index-watchdog
RuntimeDirectory=index-watchdog
RuntimeDirectoryMode=0750
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
