# scripts/index-watchdog.timer
[Unit]
Description=Check Index container health every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
